# salt.sv.conf
[program:hanlon]
command=/bin/bash -c "salt-master"
redirect_stderr=true
